# etc/services.d/portainer/run
#!/usr/bin/with-contenv bash

#exec /portainer
